# supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=./home/mongo-arm/binaries/2.1.1/bin/mongod --dbpath /home/data 
